# blocky/rootfs/etc/s6-overlay/s6-rc.d/blocky/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Blocky
# Runs Blocky
# ==============================================================================
bashio::log.info "Starting Blocky..."

exec blocky --config /etc/blocky/config.yaml
